# api/docker/laravel-websocket.conf
[program:laravel-websocket]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan websockets:serve
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=1000
numprocs=1
redirect_stderr=true
stdout_logfile=/app/storage/logs/websocket.log
stopwaitsecs=3600
